# docs/superpowers/specs/launchd/com.llm-team.dev.plist
<?xml version="1.0" encoding="UTF-8"?>
<!-- com.llm-team.dev.plist — see po.plist for documentation. -->
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
  <key>Label</key>
  <string>com.llm-team.dev</string>

  <key>ProgramArguments</key>
  <array>
    <string>__LLM_TEAM_ROOT__/scheduler/daemon.sh</string>
    <string>dev</string>
  </array>

  <key>WorkingDirectory</key>
  <string>__LLM_TEAM_ROOT__</string>

  <key>RunAtLoad</key>
  <true/>

  <key>KeepAlive</key>
  <true/>

  <key>ThrottleInterval</key>
  <integer>10</integer>

  <key>EnvironmentVariables</key>
  <dict>
    <key>PATH</key>
    <string>/opt/homebrew/bin:/usr/local/bin:/usr/bin:/bin:/usr/sbin:/sbin</string>
    <key>HOME</key>
    <string>__USER_HOME__</string>
    <key>LLM_TEAM_ROOT</key>
    <string>__LLM_TEAM_ROOT__</string>
  </dict>

  <key>StandardOutPath</key>
  <string>__LLM_TEAM_ROOT__/workdir/daemon-dev.out</string>
  <key>StandardErrorPath</key>
  <string>__LLM_TEAM_ROOT__/workdir/daemon-dev.err</string>
</dict>
</plist>
